# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:reader]
command=npm start /opt/reader
user=daemon
redirect_stderr=true
stdout_logfile=/var/opt/reader/stdout.log
stdout_logfile_maxbytes=1MB
