# etc/cron.d/clean-reader-db
# Clean old entries in DB
00 01 * * * node [ -x /opt/node/reader/current/bin/clean-db.js ] && APP_DAYS_TO_KEEP=60 /opt/node/reader/current/bin/clean-db.js
